# etc/barnacles-index.timer
[Unit]
Description=barnacles indexer dot timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target